# machines/breen/dot-config/systemd/user/fortivpn@.service
[Unit]
Description=OpenFortiVPN with SAML (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=FORTI_SERVER=%i
ExecStart=/bin/sh -c 'xdg-open-saml $FORTI_SERVER | sudo /usr/bin/openfortivpn $FORTI_SERVER --cookie-on-stdin'
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target

